# etc/oscar-website.path
[Unit]
Description="Check whether OSCAR webhook requested update"

[Path]
PathChanged=/home/oscar-www/oscar-website.trigger
Unit=oscar-website.service

[Install]
WantedBy=multi-user.target
